# templates/create-app-dirs.service.j2
[Unit]
Description=Creates application directories on startup (if they don't exist)

[Service]
ExecStart=/usr/local/bin/create-app-dirs.sh

[Install]
WantedBy=default.target
